# proprietary/intel/acd_0.2/files/service_files/com.intel.crashdump.service
[Unit]
Description=Intel BMC CPU Crashdump

[Service]
Restart=always
ExecStart=/usr/bin/crashdump
Type=dbus
BusName=com.intel.crashdump

[Install]
WantedBy=multi-user.target
